# packaging/postinstall.sh
#!/bin/sh
systemctl daemon-reload
systemctl enable lplex || true
